# etc/emailrelay.service.in
[Unit]
Description=E-MailRelay mail server
Documentation=man:emailrelay(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
Restart=on-success
KillMode=control-group
ExecStart=__SBIN_DIR__/emailrelay --as-server --syslog --pid-file __RUN_DIR__/emailrelay.pid __SYSCONF_DIR__/emailrelay.conf
ExecStop=/bin/kill -15 $MAINPID
PIDFile=__RUN_DIR__/emailrelay.pid

[Install]
WantedBy=multi-user.target

